# root/etc/services.d/guacd/run
#! /bin/bash

# Run guacd
/usr/sbin/guacd -f -b 0.0.0.0 -l 4822
